# sudo_make_me_a_report/sudo_make_me_a_report.cron
# weekly reports
30 23 * * 0	root	cd /root && /root/sudo_make_me_a_report.pl
